# telegram/telegram-bot.service
[Unit]
Description=Alfred Telegram Listener Bot
After=network.target

[Service]
Type=simple
User=sefa
WorkingDirectory=/home/sefa/alfred-hub/command-center/telegram
ExecStart=/usr/bin/python3 /home/sefa/alfred-hub/command-center/telegram/telegram_bot.py
Restart=always
RestartSec=10
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
